# systemd-examples/example_2_timer_service/systemd-cleanup-service.timer

[Unit]
Description=systemd-cleanup-service.timer

[Timer]
OnCalendar=minutely
AccuracySec=1usec

[Install]
WantedBy=multi-user.target
