# systemd-examples/example_5_nodejs_test_execution_runner/nodejs-runner-cleaner.timer

[Unit]
Description=systemd-cleanup-service.timer

[Timer]
OnCalendar=*-*-* 23:00:00
AccuracySec=1usec

[Install]
WantedBy=multi-user.target
